# malioglasi.service
[Unit]
Description=Web scraper for mobilnisvet.com

[Service]
Type=simple
ExecStart=/usr/local/bin/malioglasi.py
