# malioglasi.service
[Unit]
Description=Web scraper for http://mobilnisvet.com/mobilni-malioglasi

[Service]
Type=oneshot
ExecStart=/usr/local/bin/malioglasi.py
